# monitor/systemd/gs-rest-monitor.service
[Unit]
Description=gs-rest-service availability monitor
Documentation=https://github.com/amayabdaniel/bluegrid-devops-task
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=monitor
Group=monitor
WorkingDirectory=/opt/gs-rest-monitor
EnvironmentFile=-/etc/gs-rest-monitor.env
# Target URL and interval must be in /etc/gs-rest-monitor.env; SLACK_WEBHOOK_URL is optional.
ExecStart=/opt/gs-rest-monitor/.venv/bin/python -m gs_rest_monitor \
    --target ${TARGET_URL} \
    --interval ${INTERVAL:-15} \
    --state-file /var/lib/gs-rest-monitor/state.json
Restart=always
RestartSec=5
TimeoutStopSec=15

# --- Hardening ---
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectHostname=true
ProtectClock=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
CapabilityBoundingSet=
AmbientCapabilities=
UMask=0077
ReadWritePaths=/var/lib/gs-rest-monitor

# Logging: captured by journald, machine-parseable (already JSON)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
